# wifun/wifunc/files/wifunc.init
#!/bin/sh /etc/rc.common

START=99

start() {
	/usr/sbin/wifunc &
}

stop() {
	killall  wifunc
}
